# supervisord.conf
[supervisord]
nodaemon=false

[program:proxy]
command=./proxy.py
autorestart=true
